# codex-playwright/deploy/mcp-playwright.service.template
[Unit]
Description=codex-playwright server (SSE Transport)
After=network.target
Documentation=https://github.com/cooneycw/codex-power-pack

[Service]
Type=simple
User=${SERVICE_USER}
WorkingDirectory=${MCP_SERVER_DIR}

# Load environment variables from .env file (if exists)
EnvironmentFile=-${MCP_SERVER_DIR}/.env

# Explicit port configuration
Environment="SERVER_PORT=9101"
Environment="SERVER_HOST=127.0.0.1"

# Use uv to run in the project environment
ExecStart=${UV_BIN}/uv run --project ${MCP_SERVER_DIR} python src/server.py

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitBurst=5
StartLimitIntervalSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=codex-playwright

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=${MCP_SERVER_DIR}

[Install]
WantedBy=multi-user.target
